# deploy-sandbox/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_error.log

[program:copaw-api]
command=/app/venv/bin/copaw app --host 0.0.0.0 --port 8080
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/copaw.log
stderr_logfile=/var/log/supervisor/copaw_error.log
environment=COPAW_PORT="8080",COPAW_CONSOLE_STATIC_DIR="/usr/share/nginx/html/copaw",PYTHONUNBUFFERED="1"
